# wayland/nfs/mnt-nfs-media.automount
[Unit]
Description=Automount media

[Automount]
Where=/mnt/nfs/media

[Install]
WantedBy=multi-user.target
