# wayland/nfs/mnt-nfs-src.automount
[Unit]
Description=Automount src

[Automount]
Where=/mnt/nfs/src

[Install]
WantedBy=multi-user.target
